# app/core/unified/telegram-bot.service
[Unit]
Description=Telegram Bot Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=telegram-bot
Group=telegram-bot
WorkingDirectory=/opt/telegram-bot
ExecStart=/usr/bin/python3 /opt/telegram-bot/unified/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bot
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 